# etc/libp2p-relay-daemon.service
[Unit]
Description=lotus traced
After=network.target

[Service]
ExecStart=/root/libp2p-relay-daemon.sh
PIDFile=/var/run/libp2p-relay-daemon.pid
